# opencti-deploy-offline/config/elasticsearch.service
[Unit]
Description=Elasticsearch 8.x
Documentation=https://www.elastic.co
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# ES bắt buộc chạy bằng user riêng (hard-coded check trong bin/elasticsearch)
# Không thể chạy root — sẽ exit ngay với lỗi "can not run elasticsearch as root"
User=elasticsearch
Group=elasticsearch
RuntimeDirectory=elasticsearch
PrivateTmp=true

WorkingDirectory=/opt/elasticsearch

Environment=ES_HOME=/opt/elasticsearch
Environment=ES_PATH_CONF=/opt/elasticsearch/config
Environment=ES_JAVA_HOME=/opt/elasticsearch/jdk
Environment=PID_DIR=/var/run/elasticsearch
Environment=ES_TMPDIR=/opt/elasticsearch/tmp

ExecStart=/opt/elasticsearch/bin/elasticsearch

StandardOutput=journal
StandardError=inherit

LimitNOFILE=65535
LimitNPROC=4096
LimitAS=infinity
LimitFSIZE=infinity

TimeoutStopSec=120

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
