# tools/source-maltrail-sync/deploy/v2-ioc-maltrail-sync.service
[Unit]
Description=Maltrail IOC sync — clone + compare + update MySQL
Wants=network-online.target
After=network-online.target mysql.service

[Service]
Type=oneshot
WorkingDirectory=/opt/tools
ExecStart=/opt/tools/v2-ioc-maltrail-sync
EnvironmentFile=/etc/saids/opencti/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2-ioc-maltrail-sync

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/tools/data
PrivateTmp=yes
