# opencti-deploy-offline/systemd/opencti-platform.service
[Unit]
Description=OpenCTI Platform (GraphQL API + Frontend)
Documentation=https://docs.opencti.io
After=network-online.target redis.service minio.service rabbitmq.service elasticsearch.service
Wants=network-online.target
Requires=redis.service minio.service rabbitmq.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root

# Shared config — start script reads this .env
EnvironmentFile=/etc/saids/opencti/.env
ExecStart=/bin/bash /usr/local/bin/v2_start_opencti.sh

WorkingDirectory=/etc/saids/opencti

Restart=always
RestartSec=10

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Environment (bổ sung, start.sh override)
Environment=NODE_ENV=production

# Ensure log + data directories exist
ExecStartPre=/bin/mkdir -p /var/log/opencti
ExecStartPre=/bin/mkdir -p /etc/saids/opencti/logs
ExecStartPre=/bin/mkdir -p /etc/saids/opencti/telemetry

# Logging
StandardOutput=append:/var/log/opencti/opencti-platform.log
StandardError=append:/var/log/opencti/opencti-platform.log
SyslogIdentifier=opencti-platform

# Timeouts
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
