# opencti-deploy-offline/systemd/opencti-platform.service
[Unit]
Description=OpenCTI Platform (GraphQL API + Frontend)
Documentation=https://docs.opencti.io
After=network-online.target redis.service minio.service rabbitmq-server.service elasticsearch.service
Wants=network-online.target
Requires=redis.service minio.service rabbitmq-server.service

[Service]
Type=simple
User=root
Group=root

# Start script chứa env vars + exec npm run serv
ExecStart=/bin/bash /etc/saids/application/opencti/start.sh

WorkingDirectory=/etc/saids/application/opencti

Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Environment (bổ sung, start.sh override)
Environment=NODE_ENV=production

# Ensure log + data directories exist
ExecStartPre=/bin/mkdir -p /var/log/application/opencti
ExecStartPre=/bin/mkdir -p /etc/saids/application/opencti/logs
ExecStartPre=/bin/mkdir -p /etc/saids/application/opencti/telemetry

# Logging
StandardOutput=append:/var/log/application/opencti/opencti-platform.log
StandardError=append:/var/log/application/opencti/opencti-platform.log
SyslogIdentifier=opencti-platform

# Timeouts
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
