# v2-connectors/monitor/deploy/connector-summary.service
[Unit]
Description=OpenCTI Connector Event Summary (daily)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/tools
ExecStart=/opt/tools/connector-summary
EnvironmentFile=/etc/saids/opencti/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=connector-summary

NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
